# resources/autotest/.config/autostart/autotest.desktop
[Desktop Entry]
Type=Application
Exec=/home/xxx/autotest/autotest.sh
Name=Autotest Script
Comment=Run after desktop
