# resources/config/pgy/rc.local
#!/bin/sh
/etc/init.d/pgyvpn start &
/usr/bin/expect /xxx/pgy_login
